# etc/waterwheel-server.service
[Unit]
Description=Waterwheel Server
After=network-online.target

[Service]
Type=simple
User=ec2-user
EnvironmentFile=/etc/default/waterwheel

ExecStart=/usr/bin/waterwheel server
Restart=on-failure
RestartSec=10
KillMode=process

[Install]
WantedBy=multi-user.target
